# dot_config/systemd/user/paru-sync.timer
[Unit]
Description=Periodic paru sync

[Timer]
OnBootSec=1m
OnUnitActiveSec=60m
Persistent=true

[Install]
WantedBy=timers.target
