# dot_config/systemd/user/paru-sync.service
[Unit]
Description=Sync pacman and AUR databases

[Service]
ExecStart=/usr/bin/paru -Sy --noconfirm
