# dot_config/systemd/user/theme-dark.service
[Unit]
Description=Set Dark Theme

[Service]
Type=oneshot
ExecStart=%h/.local/bin/enable-theme-dark.sh
